# 99-auto-ripper.rules
# udev rules for automatic CD/DVD ripping
# Place this file in /etc/udev/rules.d/

# Primary rule: Trigger on optical disc media change (with media present check)
SUBSYSTEM=="block", KERNEL=="sr[0-9]*", ENV{ID_CDROM}=="1", ACTION=="change", ENV{DISK_MEDIA_CHANGE}=="1", ENV{ID_CDROM_MEDIA}=="1", RUN+="/opt/auto-ripper/trigger-rip.sh %k"

# Fallback rule for drives that don't set DISK_MEDIA_CHANGE properly
SUBSYSTEM=="block", KERNEL=="sr[0-9]*", ACTION=="change", ATTR{removable}=="1", ENV{ID_CDROM_MEDIA}=="1", RUN+="/opt/auto-ripper/trigger-rip.sh %k"